# launchd/com.maestro.discord.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.maestro.discord</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/maestro-discord</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
  </dict>

  <key>StandardOutPath</key>
  <string>/opt/homebrew/var/log/maestro-discord/output.log</string>

  <key>StandardErrorPath</key>
  <string>/opt/homebrew/var/log/maestro-discord/error.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>

  <key>WorkingDirectory</key>
  <string>/opt/homebrew/libexec/maestro-discord</string>

  <key>Umask</key>
  <integer>0022</integer>
</dict>
</plist>
